# com.xhs.autoreply.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.xhs.autoreply</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/YOUR_USERNAME/xhs_auto_reply/xhs_reply.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>API_KEY</key>
        <string>your_deepseek_api_key_here</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>your_telegram_bot_token_here</string>
        <key>TELEGRAM_CHAT_ID</key>
        <string>your_telegram_chat_id_here</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/xhs_auto_reply/logs/launchd_out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/xhs_auto_reply/logs/launchd_err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
